# cookbooks/postgresql/templates/etc/systemd/system/postgresql-9.6.service.erb
.include /lib/systemd/system/postgresql-9.6.service

[Service]
Environment=PGDATA=<%= @environment_pgdata %>
